# Docker/supervisord.conf
[supervisord]
nodaemon=true

[program:postgresql]
command=/usr/lib/postgresql/10/bin/postgres -D /var/lib/postgresql/10/main -c config_file=/etc/postgresql/10/main/postgresql.conf
user=postgres
autorestart=true

[program:httpd]
command=/bin/bash -c "exec apachectl start"
# /bin/bash -c "exec apachectl start"

#[program:sshd]
# command=/usr/sbin/sshd

[program:solr]
command=java -jar start.jar
user=solr
directory=/home/solr/solr-4.4.0/example


[program:openoffice]
command=soffice --invisible --headless "--accept=socket,host=localhost,port=2002;urp;"
user=root
directory=/tmp


